# deploy/arb.service
[Unit]
Description=Polymarket Convergence Arbitrage Bot (paper)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sniper
Group=sniper
WorkingDirectory=/opt/sniper
EnvironmentFile=/opt/sniper/.env
ExecStart=/opt/sniper/.venv/bin/python -m src.arb.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arb

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/sniper

[Install]
WantedBy=multi-user.target
